# deploy/kalima.service
[Unit]
Description=Kalima Quranic Research Platform
After=network.target

[Service]
Type=simple
User=kalima
WorkingDirectory=/opt/kalima
Environment="RUST_LOG=info"
Environment="KALIMA_DB=/opt/kalima/kalima.db"
Environment="KALIMA_INDEX=/opt/kalima/kalima-index"
ExecStart=/opt/kalima/kalima-api
Restart=on-failure

[Install]
WantedBy=multi-user.target
